# deploy/yuiolink-backup.timer
[Unit]
Description=Nightly YuioLink database snapshot

[Timer]
OnCalendar=*-*-* 03:17:00 UTC
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
